# sartor/memory/projects/family-thread-dossier/_history/pushpeer-task-pre-fix.xml
<?xml version="1.0" encoding="UTF-16"?>

<Task version="1.3" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">

  <RegistrationInfo>

    <Description>Push fresh Claude OAuth credentials to peer machines every 4 hours</Description>

    <URI>\PushPeerCredentials</URI>

  </RegistrationInfo>

  <Principals>

    <Principal id="Author">

      <UserId>S-1-5-21-2515504145-1138363155-2256348784-1001</UserId>

      <LogonType>InteractiveToken</LogonType>

    </Principal>

  </Principals>

  <Settings>

    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>

    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>

    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>

    <StartWhenAvailable>true</StartWhenAvailable>

    <IdleSettings>

      <Duration>PT10M</Duration>

      <WaitTimeout>PT1H</WaitTimeout>

      <StopOnIdleEnd>true</StopOnIdleEnd>

      <RestartOnIdle>false</RestartOnIdle>

    </IdleSettings>

    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>

  </Settings>

  <Triggers>

    <TimeTrigger>

      <StartBoundary>2026-04-29T16:31:15-04:00</StartBoundary>

      <Repetition>

        <Interval>PT1H</Interval>

        <StopAtDurationEnd>true</StopAtDurationEnd>

      </Repetition>

    </TimeTrigger>

  </Triggers>

  <Actions Context="Author">

    <Exec>

      <Command>C:\Program Files\Git\bin\bash.exe</Command>

      <Arguments>-c 'C:\Users\alto8\scripts\push-peer-credentials.sh'</Arguments>

    </Exec>

  </Actions>

</Task>
