# sartor/memory/machines/rtxpro6000server/systemd/sartor-claude-peer.service
[Unit]
Description=Sartor Claude Code peer in persistent tmux session
Documentation=file:///home/alton/Sartor-claude-network/.claude/skills/peer-comms/SKILL.md
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
WorkingDirectory=/home/alton/Sartor-claude-network
Environment=HOME=/home/alton
Environment=TERM=xterm-256color
ExecStartPre=-/usr/bin/tmux kill-session -t claude-team-1
ExecStart=/usr/bin/tmux new-session -d -s claude-team-1 -x 200 -y 50 "cd /home/alton/Sartor-claude-network && /usr/bin/claude --dangerously-skip-permissions"
ExecStop=/usr/bin/tmux kill-session -t claude-team-1
RemainAfterExit=yes
Restart=on-failure
RestartSec=15

[Install]
WantedBy=default.target
